# avarc-chatopts-bot-restart.service
[Unit]
Description=Restart avarc-chatopts-bot Service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart avarc-chatopts-bot.service
